# home/private_dot_config/systemd/user/restic-backup.timer
[Unit]
Description=Hourly restic local backup

[Timer]
OnCalendar=hourly
Persistent=true
RandomizedDelaySec=300

[Install]
WantedBy=timers.target
